# deploy/systemd/gameops-bot.service
[Unit]
Description=GameOps Bridge Discord Bot
After=network.target

[Service]
Type=simple
WorkingDirectory=/Users/christopherrossi/Desktop/Projects/gameops-bridge
Environment=NODE_ENV=production
Environment=GAMEOPS_CONFIG_PATH=/tmp/gameops.bootstrap.sample.json
Environment=BOT_LOCAL_CONFIG_PATH=/Users/christopherrossi/Desktop/Projects/gameops-bridge/config/bot.local.json
ExecStart=/usr/bin/env npm --workspace apps/bot run dev
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
